# deploy/systemd/video-web.service
[Unit]
Description=Video Generator Flask Web
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/video-generator
Environment="PATH=/opt/video-generator/.venv/bin"
ExecStart=/opt/video-generator/.venv/bin/gunicorn -w 2 -b 127.0.0.1:5000 app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
